# yaw-server/services/icantbelieveitsnotCI.timer
[Unit]
Description=I can't believe its not CI timer

[Timer]
OnCalendar=*:*:0/10
