# yaw-server/services/yaw-api.service
[Unit]
Description=Start yaw Express API

[Service]
ExecStart=/usr/bin/node /srv/scripts/express/app.js
