# yaw-server/services/yaw-import.timer
[Unit]
Description=Timer for mongodbimport for aircraft data

[Timer]
OnCalendar=*:*:0/10
