# server/services/yaw-socket.service
[Unit]
Description=start socket server for yaw

[Service]
ExecStart=/usr/bin/node /srv/scripts/socket/server.js